# CTF_CUH/18_linux_expuesto_servicio_olvidado/solutions/config/cuh-reportes.service
[Unit]
Description=CUH reportes legado

[Service]
User=cuh-reportes
ExecStart=/usr/bin/python3 /opt/cuh/reportes/server.py --host 127.0.0.1
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
